# inittab
::respawn:/usr/local/bin/start_syslogd
::respawn:/usr/local/bin/start_prometheus
::respawn:/usr/local/bin/start_grafana
